# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:router]
command=python /app/run_demo.py
directory=/app
autostart=false
autorestart=false
stdout_logfile=/var/log/supervisor/router.log
stderr_logfile=/var/log/supervisor/router_err.log

[program:pod-0]
command=python /app/workers/model_pod.py --pod-id=0
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pod-0.log
stderr_logfile=/var/log/supervisor/pod-0_err.log

[program:pod-1]
command=python /app/workers/model_pod.py --pod-id=1
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pod-1.log
stderr_logfile=/var/log/supervisor/pod-1_err.log

[program:pod-2]
command=python /app/workers/model_pod.py --pod-id=2
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pod-2.log
stderr_logfile=/var/log/supervisor/pod-2_err.log

[program:pod-3]
command=python /app/workers/model_pod.py --pod-id=3
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pod-3.log
stderr_logfile=/var/log/supervisor/pod-3_err.log

[program:pod-4]
command=python /app/workers/model_pod.py --pod-id=4
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pod-4.log
stderr_logfile=/var/log/supervisor/pod-4_err.log
